# debian/wikigps-tracker.postinst
#!/bin/bash
cd /etc/init.d;
update-rc.d defaults 99;
chmod a+w /var/wikigps-tracker;
