# debian/wikigps-tracker.postinst
#!/bin/sh
/usr/bin/wikigpsTracker > /dev/null 2>&1 &
update-rc.d tracker defaults;
exit 0
