# debian/wikigps-tracker.postinst
#!/bin/sh
/usr/bin/wikigpsTracker&
update-rc.d tracker defaults;
exit 0
