# debian/wikigps-tracker.postinst
#!/bin/bash
#cd /etc/init.d;
touch /var/wikigps-tracker;
chmod a+w /var/wikigps-tracker;
update-rc.d tracker defaults;
